# scripts/backup/datafoundry-backup.timer
[Unit]
Description=Run DataFoundry Postgres backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
